# systemd/whisper-hotkey.service.template
[Unit]
Description=Whisper Dictation Hotkey Listener
After=graphical.target whisper-daemon.service
Wants=whisper-daemon.service
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
ExecStart=/usr/bin/python3 $HOME/.local/bin/whisper-hotkey
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
